# hotel_api/supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock                   ; (the path to the socket file)
chmod=0777                                      ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/carlypso/ford_quic_supervisor.log           ; supervisord log file
logfile_maxbytes=1MB                           ; maximum size of logfile before rotation
logfile_backups=0                              ; number of backed up logfiles
loglevel=debug                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/carlypso/                  ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:apache]
command=python3 manage.py runserver
autostart=true
autorestart=true
startretries=1
startsecs=1
redirect_stderr=true
stderr_logfile=/var/log/hotel_api/mylog.err.log
stdout_logfile=/var/log/hotel_api/mylog.out.log
user=root
killasgroup=true
stopasgroup=true
